# helperscripts/unattended-upgrade/my-unattended-upgrade.cron
# m h dom mon dow user  command
26 * * * * root /usr/local/bin/my-unattended-upgrade
