# helperscripts/scheduled-restart/scheduled-restart.cron
20 6    * * *   root   /usr/local/bin/scheduled-restart
