# helperscripts/scheduled-restart/scheduled-restart.cron
@reboot	root	/bin/rm /var/run/need-system-restart
20 6	* * *	root	/bin/sh /usr/local/bin/scheduled-restart
