# helperscripts/autoupdate/autoupdate.cron
8 *	* * *	root	/usr/local/bin/autoupdate
