# helperscripts/autoupdate/cron/autoupdate.cron
8 *	* * *	root	/bin/sh /usr/local/bin/autoupdate
