# helperscripts/autoupdate/systemd/systemd/autoupdate.timer
[Unit]
Description=Autoupdate timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=autoupdate.target

[Install]
WantedBy=basic.target
