# helperscripts/autoupdate/systemd/systemd/autoupdate.target
[Unit]
Description=Autoupdate target
Wants=autoupdate.service

StopWhenUnneeded=yes
